# conf/gowebserver.conf
[program:goserver]
command=/usr/local/bin/simple-go-webserver
autostart=true
autorestart=true
startretries=3
killasgroup=true
stopasgroup=true
